# data_logger.service
[Unit]
Description=Data Logger Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/data-logger/data_logger.py
Restart=always
RestartSec=3
User=root
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target